# sync/com.primebot.knowledgesync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.primebot.knowledgesync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/anasahmadouch/Desktop/primebot-knowledge/sync/sync_to_obsidian.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/anasahmadouch/Desktop/primebot-knowledge/sync/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/anasahmadouch/Desktop/primebot-knowledge/sync/sync.error.log</string>
</dict>
</plist>
